# source/systemd/chip-battery-status-update.service
[Unit]
Description=Get C.H.I.P. battery status

[Service]
ExecStart=/usr/local/sbin/chip-battery-update-status